# debian/cron.daily
#!/bin/sh
/usr/sbin/mozilla-nightly-updater
